# scripts/services/mnt-volume.automount
[Unit]
Description=Automount usb

[Automount]
Where=/mnt/volume

[Install]
WantedBy=multi-user.target
